# deployment/systemd/daily-update.service
[Unit]
Description=NBA Simulator - Daily ESPN Update Workflow
Documentation=https://github.com/ryanranft/nba-simulator-aws
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ryanranft
Group=staff
WorkingDirectory=/Users/ryanranft/nba-simulator-aws
Environment="PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/Users/ryanranft/nba-simulator-aws"

# Execute workflow
ExecStart=/opt/homebrew/bin/python3 scripts/workflows/daily_update_cli.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nba-daily-update

# Resource limits
MemoryMax=1G
CPUQuota=100%
TimeoutStartSec=900
TimeoutStopSec=180

# Restart policy
Restart=on-failure
RestartSec=120
StartLimitBurst=5
StartLimitIntervalSec=24h

[Install]
WantedBy=multi-user.target
