# deployment/systemd/validation.service
[Unit]
Description=NBA Simulator - 3-Source Validation Workflow
Documentation=https://github.com/ryanranft/nba-simulator-aws
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ryanranft
Group=staff
WorkingDirectory=/Users/ryanranft/nba-simulator-aws
Environment="PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/Users/ryanranft/nba-simulator-aws"

# Execute workflow
ExecStart=/opt/homebrew/bin/python3 scripts/workflows/validation_cli.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nba-validation

# Resource limits
MemoryMax=2G
CPUQuota=150%
TimeoutStartSec=1800
TimeoutStopSec=300

# Restart policy
Restart=on-failure
RestartSec=180
StartLimitBurst=3
StartLimitIntervalSec=24h

[Install]
WantedBy=multi-user.target
